# conf/supervise-ngircd.conf
[program:Supervise_ngircd]
user=irc
command=/usr/sbin/ngircd -n -f /etc/ngircd/ngircd.conf
redirect_stderr=true
autorestart=true
