# conf/supervise-tor.conf
[program:Supervise_tor]
command=/usr/bin/tor -f /etc/tor/torrc
redirect_stderr=true
autorestart=true
